# machine/files/devops-infra/devops-platform-restart.path
[Path]
PathModified=/opt/platform/git-keys/id_rsa.pub

[Install]
WantedBy=default.target